# scripts/trader_tools/com.quantrabbit.disk-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.quantrabbit.disk-cleanup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tossaki/App/QuantRabbit/scripts/trader_tools/disk_cleanup.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/tossaki/App/QuantRabbit</string>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardErrorPath</key>
    <string>/Users/tossaki/App/QuantRabbit/logs/disk_cleanup_err.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
