# archive/systemd/quant-ssh-watchdog.service
[Unit]
Description=QuantRabbit SSH watchdog
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /home/tossaki/QuantRabbit/scripts/ssh_watchdog.sh
