# archive/systemd/quant-boot-sync.service
[Unit]
Description=QuantRabbit boot sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/qr-boot-sync

[Install]
WantedBy=multi-user.target
