# archive/systemd/quant-online-tuner.service
[Unit]
Description=QuantRabbit Online Tuner (exit/gate/alloc)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/tossaki/QuantRabbit
EnvironmentFile=-/home/tossaki/QuantRabbit/ops/env/quant-v2-runtime.env
Environment="PYTHONUNBUFFERED=1"
ExecStart=/bin/bash /home/tossaki/QuantRabbit/scripts/run_online_tuner_live.sh
User=tossaki
