# systemd/cleanup-qr-logs.timer
[Unit]
Description=Twice-daily QuantRabbit log cleanup (roughly UTC 07:30/19:30)

[Timer]
OnCalendar=*-*-* 07:30:00
OnCalendar=*-*-* 19:30:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
