# systemd/qr-macro-donchian55.service
[Unit]
Description=QR Macro Donchian55 Worker (single strategy)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
Environment="PYTHONUNBUFFERED=1"
Environment="MACRO_STRATEGY_ALLOWLIST=Donchian55"
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.macro_core.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
