# systemd/qr-scalp-pulse-break.service
[Unit]
Description=QR Scalp PulseBreak Worker (single strategy)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
Environment="PYTHONUNBUFFERED=1"
Environment="SCALP_STRATEGY_ALLOWLIST=PulseBreak"
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_core.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
