# systemd/quant-scalp-precision-false-break-fade.service
[Unit]
Description=QuantRabbit Scalp Precision False Break Fade Worker (FBF-123 tuned)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
EnvironmentFile=/etc/quantrabbit.env
Environment=PYTHONUNBUFFERED=1
Environment=SCALP_PRECISION_ENABLED=1
Environment=SCALP_PRECISION_MODE=false_break_fade
Environment=SCALP_PRECISION_UNIT_ALLOWLIST=false_break_fade
Environment=SCALP_PRECISION_MODE_FILTER_ALLOWLIST=1
Environment=SCALP_PRECISION_LOOP_INTERVAL_SEC=2.0
Environment=SCALP_PRECISION_COOLDOWN_SEC=12
Environment=SCALP_PRECISION_MAX_OPEN_TRADES=1
Environment=SCALP_PRECISION_OPEN_TRADES_SCOPE=tag
Environment=SCALP_PRECISION_LOG_PREFIX=[ScalpPrecision:FalseBreakFade123]

# Tune (week replay: 2026-02-02..2026-02-06): require deeper sweep + reclaim to avoid churn.
Environment=FBF_BREAKOUT_PIPS=1.4
Environment=FBF_RECLAIM_PIPS=0.6
Environment=FBF_MIN_SWEEP_PIPS=1.6
Environment=FBF_ADX_MAX=26.0
Environment=FBF_RANGE_SCORE_MIN=0.45
Environment=FBF_ATR_MAX=4.2
Environment=FBF_MAX_COUNTER_SLOPE_PIPS=0.45
Environment=FBF_MAX_COUNTER_VWAP_GAP_PIPS=45.0
Environment=FBF_TICK_MIN_STRENGTH=0.35
Environment=FBF_LOOKBACK=30
Environment=FBF_HI_PCT=97.0
Environment=FBF_LO_PCT=3.0

# Broker SL/TP disabled (no entry stop-loss). Re-enable by setting:
# Environment=ORDER_ENABLE_STOP_LOSS_SCALP=1

ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_precision.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
