# systemd/quant-onepip-s1.service
[Unit]
Description=QuantRabbit OnePip Maker S1 Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.onepip_maker_s1.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
