# systemd/quant-scalp-squeeze-pulse-break.service
[Unit]
Description=QuantRabbit Scalp Squeeze Pulse Break Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
EnvironmentFile=/home/tossaki/QuantRabbit/ops/env/quant-v2-runtime.env
Environment=PYTHONUNBUFFERED=1
Environment=SCALP_PRECISION_ENABLED=1
Environment=SCALP_PRECISION_MODE=squeeze_pulse_break
Environment=SCALP_PRECISION_UNIT_ALLOWLIST=squeeze_pulse_break
Environment=SCALP_PRECISION_MODE_FILTER_ALLOWLIST=1
Environment=SCALP_PRECISION_LOOP_INTERVAL_SEC=2.0
Environment=SCALP_PRECISION_COOLDOWN_SEC=10
Environment=SCALP_PRECISION_MAX_OPEN_TRADES=1
Environment=SCALP_PRECISION_OPEN_TRADES_SCOPE=tag
Environment=SPB_BBW_MAX=0.0020
Environment=SPB_ATR_MIN=0.9
Environment=SPB_ATR_MAX=4.6
Environment=SPB_LOOKBACK=18
Environment=SPB_HI_PCT=85.0
Environment=SPB_LO_PCT=15.0
Environment=SPB_BREAKOUT_PIPS=0.6
Environment=SPB_RANGE_SCORE_MIN=0.58
Environment=SPB_ALLOWED_REGIMES=trend
Environment=SPB_TICK_WINDOW_SEC=12.0
Environment=SPB_MIN_TICKS=10
Environment=SPB_TICK_RANGE_MAX_PIPS=1.8
Environment=SPB_IMB_RATIO_MIN=0.60
Environment=SPB_MOM_MIN_PIPS=0.55
Environment=SPB_SPREAD_P25=1.0
Environment=SPB_SIZE_MULT=1.15
Environment=SPB_REQUIRE_AIR_MATCH=1
Environment=SPB_MIN_AIR_SCORE=0.60
Environment=SPB_ALLOW_LONG=0
Environment=SPB_ALLOW_SHORT=1
Environment=SPB_MAX_ADX=30
Environment=SPB_SHORT_RSI_MIN=20
Environment=SPB_SHORT_RSI_MAX=55
Environment=SPB_BLOCK_JST_HOURS=18
Environment=SPB_DIAG=1
Environment=SPB_DIAG_INTERVAL_SEC=15
Environment=ORDER_ENABLE_STOP_LOSS_SCALP=1
Environment=ORDER_ALLOW_STOP_LOSS_WITH_EXIT_NO_NEGATIVE_CLOSE_SCALP=true
Environment=ORDER_ENTRY_HARD_STOP_PIPS_STRATEGY_SQUEEZEPULSEBREAK=7.0
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_squeeze_pulse_break.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
