# systemd/quant-scalp-precision-rangefaderpro-exit.service
[Unit]
Description=QuantRabbit Scalp Precision RangeFaderPro Exit Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
EnvironmentFile=/etc/quantrabbit.env
Environment=PYTHONUNBUFFERED=1
Environment=SCALP_PRECISION_EXIT_TAGS=RangeFaderPro
Environment=SCALP_PRECISION_EXIT_LOG_PREFIX=[ScalpPrecisionExit:RangeFaderPro]
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_precision.exit_worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
