# systemd/quant-scalp-precision-tick-imbalance.service
[Unit]
Description=QuantRabbit Scalp Precision Tick Imbalance Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
EnvironmentFile=/etc/quantrabbit.env
Environment=PYTHONUNBUFFERED=1
Environment=SCALP_PRECISION_ENABLED=1
Environment=SCALP_PRECISION_MODE=tick_imbalance
Environment=SCALP_PRECISION_ALLOWLIST=tick_imbalance
Environment=SCALP_PRECISION_MODE_FILTER_ALLOWLIST=1
Environment=SCALP_PRECISION_MAX_OPEN_TRADES=1
Environment=SCALP_PRECISION_COOLDOWN_SEC=180
Environment=SCALP_PRECISION_LOG_PREFIX=[ScalpPrecision:TickImb]
Environment=TICK_IMB_RATIO_MIN=0.74
Environment=TICK_IMB_MOM_MIN_PIPS=0.65
Environment=TICK_IMB_RANGE_MIN_PIPS=0.35
Environment=TICK_IMB_ADX_MIN=22.0
Environment=TICK_IMB_MA_GAP_MIN_PIPS=0.12

# Broker SL/TP disabled (no entry stop-loss).  Re-enable by setting:
# Environment=ORDER_ENABLE_STOP_LOSS_SCALP=1

ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_precision.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
