# systemd/quant-scalp-precision-liquidity-sweep-hf.service
[Unit]
Description=QuantRabbit Scalp Precision LiquiditySweep Reversal HF Worker (LSR-126 tuned)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
EnvironmentFile=/etc/quantrabbit.env
Environment=PYTHONUNBUFFERED=1
Environment=SCALP_PRECISION_ENABLED=1
Environment=SCALP_PRECISION_MODE=liquidity_sweep
Environment=SCALP_PRECISION_UNIT_ALLOWLIST=liquidity_sweep
Environment=SCALP_PRECISION_MODE_FILTER_ALLOWLIST=1
Environment=SCALP_PRECISION_LOOP_INTERVAL_SEC=2.0
Environment=SCALP_PRECISION_COOLDOWN_SEC=12
Environment=SCALP_PRECISION_MAX_OPEN_TRADES=1
Environment=SCALP_PRECISION_OPEN_TRADES_SCOPE=tag
Environment=SCALP_PRECISION_LOG_PREFIX=[ScalpPrecision:LiquiditySweep126]

# HF tuning (week replay: 2026-02-02..2026-02-06): stricter sweep + reclaim to cut false positives.
Environment=LSR_SWEEP_PIPS=0.65
Environment=LSR_RECLAIM_PIPS=0.25
Environment=LSR_ADX_MAX=26.0
Environment=LSR_RANGE_SCORE_MIN=0.40
Environment=LSR_WICK_RATIO_MIN=0.55

# Broker SL disabled (no entry stop-loss). Re-enable by setting:
# Environment=ORDER_ENABLE_STOP_LOSS_SCALP=1

ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_precision.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target

