# systemd/quant-autotune-ui.service
[Unit]
Description=QuantRabbit Autotune Review UI
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/uvicorn apps.autotune_ui:app --host 0.0.0.0 --port 8088
Restart=on-failure
RestartSec=5
User=tossaki
Environment=PYTHONUNBUFFERED=1
Environment=UI_SNAPSHOT_LITE_MODE=full
EnvironmentFile=-/etc/quantrabbit.env
Environment=AUTOTUNE_BQ_TABLE=quantrabbit.autotune_runs
Environment=AUTOTUNE_BQ_SETTINGS_TABLE=quantrabbit.autotune_settings

[Install]
WantedBy=multi-user.target
