# systemd/quant-realtime-metrics.service
[Unit]
Description=QuantRabbit Realtime Metrics (BigQuery)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/tossaki/QuantRabbit
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python analytics/realtime_metrics_job.py --log-level INFO
User=tossaki
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/quantrabbit.env
