# systemd/quant-trend-reclaim-long.service
[Unit]
Description=QuantRabbit Trend Reclaim Long Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
EnvironmentFile=-/etc/quantrabbit.env
Environment="PYTHONUNBUFFERED=1"
Environment="FORECAST_GATE_ENABLED=0"
Environment="ORDER_ENTRY_QUALITY_GATE_ENABLED=1"
Environment="ORDER_ENTRY_QUALITY_MICROSTRUCTURE_ENABLED=1"
Environment="TRRECLAIM_ENABLED=1"
Environment="TRRECLAIM_LOOP_INTERVAL_SEC=5.0"
Environment="TRRECLAIM_COOLDOWN_SEC=45"
Environment="TRRECLAIM_BASE_UNITS=6500"
Environment="TRRECLAIM_MIN_ENTRY_CONF=56"
Environment="TRRECLAIM_MAX_SPREAD_PIPS=1.2"
Environment="TRRECLAIM_DYN_ALLOC_ENABLED=1"
Environment="TRRECLAIM_DYN_ALLOC_LOSER_BLOCK=1"
Environment="TRRECLAIM_DYN_ALLOC_MIN_TRADES=8"
Environment="TRRECLAIM_DYN_ALLOC_LOSER_SCORE=0.30"
Environment="TRRECLAIM_BLOCK_IF_LONG_OPEN=0"
Environment="ORDER_ENABLE_STOP_LOSS_SCALP=1"
Environment="ORDER_ALLOW_STOP_LOSS_WITH_EXIT_NO_NEGATIVE_CLOSE_SCALP=true"
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_trend_reclaim.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
