# systemd/quant-scalp-tick-imbalance.service
[Unit]
Description=QuantRabbit Scalp Tick Imbalance Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/quantrabbit.env
Environment=SCALP_PRECISION_ENABLED=1
Environment=SCALP_PRECISION_MODE=tick_imbalance
Environment=SCALP_PRECISION_ALLOWLIST=tick_imbalance
Environment=SCALP_PRECISION_MODE_FILTER_ALLOWLIST=1
Environment=SCALP_PRECISION_MAX_OPEN_TRADES=1
Environment=SCALP_PRECISION_COOLDOWN_SEC=180
Environment=SCALP_PRECISION_LOG_PREFIX=[Scalp:TickImb]
Environment=SCALP_PRECISION_UNIT_PROFIT_GUARD_SCOPE=strategy
Environment=TICK_IMB_RATIO_MIN=0.74
Environment=TICK_IMB_MOM_MIN_PIPS=0.65
Environment=TICK_IMB_RANGE_MIN_PIPS=0.35
Environment=TICK_IMB_ADX_MIN=22.0
Environment=TICK_IMB_MA_GAP_MIN_PIPS=0.12
Environment=TICK_IMB_ALLOWED_REGIMES=trend
Environment=TICK_IMB_REENTRY_LOOKBACK_SEC=420
Environment=TICK_IMB_REENTRY_MIN_PRICE_GAP_PIPS=1.0
Environment=TICK_IMB_REENTRY_REQUIRE_LAST_PROFIT=1
Environment=ORDER_ENABLE_STOP_LOSS_SCALP=1
Environment=ORDER_ALLOW_STOP_LOSS_WITH_EXIT_NO_NEGATIVE_CLOSE_SCALP=true
Environment=ORDER_ENTRY_HARD_STOP_PIPS_STRATEGY_TICKIMBALANCE=25
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_tick_imbalance.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
