# systemd/quant-scalp-macd-rsi-div.service
[Unit]
Description=QuantRabbit MACD RSI Divergence Scalp Worker (M1)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
EnvironmentFile=-/etc/quantrabbit.env
Environment="PYTHONUNBUFFERED=1"
Environment="MACDRSIDIV_ENABLED=1"
Environment="MACDRSIDIV_POCKET=scalp"
Environment="MACDRSIDIV_STRATEGY_TAG=scalp_macd_rsi_div_live"
Environment="MACDRSIDIV_LOOP_INTERVAL_SEC=2.0"
Environment="MACDRSIDIV_COOLDOWN_SEC=75"
Environment="MACDRSIDIV_MAX_SPREAD_PIPS=1.1"
Environment="MACDRSIDIV_MAX_ADX=40"
Environment="MACDRSIDIV_REQUIRE_RANGE_ACTIVE=0"
Environment="MACDRSIDIV_RANGE_MIN_SCORE=0.10"
Environment="MACDRSIDIV_RSI_LONG_ARM=32"
Environment="MACDRSIDIV_RSI_SHORT_ARM=68"
Environment="MACDRSIDIV_RSI_LONG_ENTRY=36"
Environment="MACDRSIDIV_RSI_SHORT_ENTRY=64"
Environment="MACDRSIDIV_MIN_DIV_SCORE=0.08"
Environment="MACDRSIDIV_MIN_DIV_STRENGTH=0.12"
Environment="MACDRSIDIV_MAX_DIV_AGE_BARS=14"
Environment="MACDRSIDIV_ALLOW_HIDDEN_DIVERGENCE=1"
Environment="MACDRSIDIV_BASE_ENTRY_UNITS=7000"
Environment="MACDRSIDIV_MIN_UNITS=1200"
Environment="MACDRSIDIV_MAX_OPEN_TRADES=1"
Environment="MACDRSIDIV_CAP_MIN=0.10"
Environment="MACDRSIDIV_CAP_MAX=0.85"
Environment="SPREAD_GUARD_DISABLE=1"
Environment="MACDRSIDIV_SL_ATR_MULT=0.85"
Environment="MACDRSIDIV_TP_ATR_MULT=1.10"
Environment="MACDRSIDIV_PATTERN_GATE_OPT_IN=1"
Environment="ORDER_ENABLE_STOP_LOSS_SCALP=1"
Environment="ORDER_ALLOW_STOP_LOSS_WITH_EXIT_NO_NEGATIVE_CLOSE_SCALP=true"
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_macd_rsi_div.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
