# systemd/quant-macro-exit.service
[Unit]
Description=QuantRabbit Macro Exit Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.macro_exit.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
