# systemd/quant-health-snapshot.timer
[Unit]
Description=QuantRabbit health snapshot timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=20s
Unit=quant-health-snapshot.service

[Install]
WantedBy=timers.target
