# systemd/quant-scalp-precision-tick-imbalance-exit.service
[Unit]
Description=QuantRabbit Scalp Precision TickImbalance Exit Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
EnvironmentFile=/etc/quantrabbit.env
Environment=PYTHONUNBUFFERED=1
Environment=SCALP_PRECISION_EXIT_TAGS=TickImbalance
Environment=SCALP_PRECISION_EXIT_LOG_PREFIX=[ScalpPrecisionExit:TickImbalance]
Environment=RANGEFADER_EXIT_LOSS_CUT_ENABLED=true
Environment=RANGEFADER_EXIT_LOSS_CUT_REQUIRE_SL=false
Environment=RANGEFADER_EXIT_LOSS_CUT_SOFT_PIPS=3.0
Environment=RANGEFADER_EXIT_LOSS_CUT_HARD_PIPS=4.5
Environment=RANGEFADER_EXIT_LOSS_CUT_MAX_HOLD_SEC=600
Environment=RANGEFADER_EXIT_LOSS_CUT_COOLDOWN_SEC=8
Environment=RANGEFADER_EXIT_LOSS_CUT_REASON_SOFT=m1_structure_break
Environment=RANGEFADER_EXIT_LOSS_CUT_REASON_HARD=max_adverse
Environment=RANGEFADER_EXIT_LOSS_CUT_REASON_TIME=time_stop
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_precision.exit_worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
