# systemd/quant-scalp-precision-precision-lowvol.service
[Unit]
Description=QuantRabbit Scalp Precision LowVol Precision Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
EnvironmentFile=/etc/quantrabbit.env
Environment=PYTHONUNBUFFERED=1
Environment=SCALP_PRECISION_ENABLED=1
Environment=SCALP_PRECISION_MODE=precision_lowvol
Environment=SCALP_PRECISION_ALLOWLIST=precision_lowvol
Environment=SCALP_PRECISION_MIN_ENTRY_CONF=45
Environment=SCALP_PRECISION_CONF_FLOOR=45
Environment=SCALP_PRECISION_CONF_CEIL=92
Environment=SCALP_PRECISION_MAX_SIGNALS_PER_CYCLE=4
Environment=SCALP_PRECISION_COOLDOWN_SEC=8
Environment=SCALP_PRECISION_MAX_OPEN_TRADES=3
Environment=SCALP_PRECISION_MAX_OPEN_TRADES_GLOBAL=4
Environment=SCALP_PRECISION_MAX_SPREAD_PIPS=1.8
Environment=SCALP_PRECISION_PREC_LOWVOL_SPREAD_P25=1.8
Environment=SCALP_PRECISION_PREC_LOWVOL_REV_MIN_STRENGTH=0.28
Environment=SCALP_PRECISION_BASE_UNITS=18000
Environment=SCALP_PRECISION_MAX_MARGIN_USAGE=0.96
Environment=SCALP_PRECISION_CAP_MAX=1.0
Environment=DYN_SIZE_BASE_RISK_PCT=0.018
Environment=DYN_SIZE_MIN_RISK_PCT=0.004
Environment=DYN_SIZE_MAX_RISK_PCT=0.05
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_precision.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
