# systemd/quant-scalp-wick-reversal-pro-exit.service
[Unit]
Description=QuantRabbit Scalp Wick Reversal Pro Exit Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
EnvironmentFile=/etc/quantrabbit.env
Environment=PYTHONUNBUFFERED=1
Environment=SCALP_PRECISION_EXIT_TAGS=WickReversalPro
Environment=SCALP_PRECISION_EXIT_LOG_PREFIX=[ScalpExit:WickRevPro]
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_wick_reversal_pro.exit_worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
