# systemd/quant-micro-multi.service
[Unit]
Description=QuantRabbit Micro Multi-Strategy Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
Environment="PYTHONUNBUFFERED=1"
Environment="FORECAST_GATE_ENABLED=0"
Environment="MICRO_STRATEGY_ALLOWLIST=MicroLevelReactor,MicroRangeBreak,MicroVWAPRevert,MomentumPulse"
Environment="MICRO_MULTI_LOOP_INTERVAL_SEC=6.0"
Environment="MICRO_MULTI_BASE_UNITS=28000"
Environment="MICRO_MULTI_MAX_MARGIN_USAGE=0.92"
Environment="MICRO_MULTI_MAX_SIGNALS_PER_CYCLE=3"
Environment="MICRO_MULTI_MULTI_SIGNAL_MIN_SCALE=0.70"
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.micro_multistrat.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
