# systemd/quant-scalp-precision-vwap-revert-exit.service
[Unit]
Description=QuantRabbit Scalp Precision VwapRevertS Exit Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
Environment=PYTHONUNBUFFERED=1
Environment=SCALP_PRECISION_EXIT_TAGS=VwapRevertS
Environment=SCALP_PRECISION_EXIT_LOG_PREFIX=[ScalpPrecisionExit:VwapRevertS]
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_precision.exit_worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
