# systemd/quant-m1scalper.service
[Unit]
Description=QuantRabbit M1Scalper Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
Environment="PYTHONUNBUFFERED=1"
Environment=M1SCALP_ALLOWED_REGIMES=range,trend,breakout
EnvironmentFile=-/etc/quantrabbit.env
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_m1scalper.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
