# systemd/quant-squeeze-break-s5.service
[Unit]
Description=QuantRabbit Squeeze Break S5 Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.squeeze_break_s5.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
