# systemd/quant-fast-scalp-rangeonly.service
[Unit]
Description=QuantRabbit FastScalp RangeOnly Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
Environment=PYTHONUNBUFFERED=1
Environment=FAST_SCALP_ENABLED=1
Environment=FAST_SCALP_RANGE_ONLY=1
Environment=FAST_SCALP_STRATEGY_TAG=fast_scalp_range
Environment=FAST_SCALP_PATTERN_ALLOWLIST=flat_chop,neutral,narrow_range,pullback_after_impulse_up,pullback_after_impulse_down
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.fast_scalp.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
