# systemd/quant-replay-quality-gate.timer
[Unit]
Description=Run quant-replay-quality-gate.service every 3 hours

[Timer]
OnBootSec=12min
OnUnitActiveSec=3h
AccuracySec=2min
RandomizedDelaySec=5min
Unit=quant-replay-quality-gate.service

[Install]
WantedBy=timers.target
