# systemd/quant-scalp-precision-tick-imbalance.service
[Unit]
Description=QuantRabbit Scalp Precision Tick Imbalance Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
EnvironmentFile=/etc/quantrabbit.env
Environment=PYTHONUNBUFFERED=1
Environment=SCALP_PRECISION_ENABLED=1
Environment=SCALP_PRECISION_MODE=tick_imbalance
Environment=SCALP_PRECISION_ALLOWLIST=tick_imbalance
Environment=SCALP_PRECISION_MODE_FILTER_ALLOWLIST=1
Environment=SCALP_PRECISION_MAX_OPEN_TRADES=1
Environment=SCALP_PRECISION_LOG_PREFIX=[ScalpPrecision:TickImb]

# Attach broker SL/TP on entry for this worker only (caps tail losses).
Environment=ORDER_ENABLE_STOP_LOSS_SCALP=1

ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_precision.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
