# systemd/quant-autotune.service
[Unit]
Description=QuantRabbit Continuous Backtest + AutoTune
After=network-online.target

[Service]
Type=simple
# ここをあなたのリポの絶対パスに変更してください
WorkingDirectory=/opt/quantrabbit
ExecStart=/usr/bin/env python3 scripts/continuous_backtest.py --trials 40 --strategies M1Scalper,PulseBreak,RangeFader
Restart=on-failure
RestartSec=10
User=%i
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
