# systemd/quant-h1momentum.service
[Unit]
Description=QuantRabbit H1Momentum Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python - <<'PY'
import asyncio
from workers.macro_h1momentum import h1momentum_worker

asyncio.run(h1momentum_worker())
PY
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
