# systemd/quant-scalp-precision-tick-wick-reversal.service
[Unit]
Description=QuantRabbit Scalp Precision Tick Wick Reversal Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
EnvironmentFile=/etc/quantrabbit.env
Environment=PYTHONUNBUFFERED=1
Environment=SCALP_PRECISION_ENABLED=1
Environment=SCALP_PRECISION_MODE=tick_wick_reversal
Environment=SCALP_PRECISION_LOOP_INTERVAL_SEC=3.0
Environment=SCALP_PRECISION_COOLDOWN_SEC=12
Environment=SCALP_PRECISION_MAX_OPEN_TRADES=1
Environment=SCALP_PRECISION_OPEN_TRADES_SCOPE=tag
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_precision.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
