# systemd/quant-strategy-optimizer.service
[Unit]
Description=QuantRabbit Strategy Optimizer (BigQuery ML)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/tossaki/QuantRabbit
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python analytics/strategy_optimizer_job.py --log-level INFO
User=tossaki
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/quantrabbit.env
