# systemd/quant-trade-counterfactual.timer
[Unit]
Description=Run quant-trade-counterfactual.service every 20 minutes

[Timer]
OnBootSec=15min
OnUnitActiveSec=20min
AccuracySec=2min
RandomizedDelaySec=2min
Unit=quant-trade-counterfactual.service

[Install]
WantedBy=timers.target
