# systemd/quant-scalp-precision-level-reject.service
[Unit]
Description=QuantRabbit Scalp Precision Level Reject Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tossaki/QuantRabbit
EnvironmentFile=/etc/quantrabbit.env
Environment=PYTHONUNBUFFERED=1
Environment=SCALP_PRECISION_ENABLED=1
Environment=SCALP_PRECISION_MODE=level_reject
Environment=SCALP_PRECISION_ALLOWLIST=level_reject
Environment=SCALP_PRECISION_MODE_FILTER_ALLOWLIST=1
Environment=SCALP_PRECISION_COOLDOWN_SEC=30
Environment=LEVEL_REJECT_SIZE_MULT=1.6
Environment=LEVEL_REJECT_BAND_PIPS=0.7
Environment=LEVEL_REJECT_ATR_MAX=3.6
Environment=LEVEL_REJECT_SPREAD_P25=1.0
ExecStart=/home/tossaki/QuantRabbit/.venv/bin/python -m workers.scalp_precision.worker
Restart=always
RestartSec=5
User=tossaki

[Install]
WantedBy=multi-user.target
