# scripts/guardian/com.quantrabbit.guardian-wake-dispatcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.quantrabbit.guardian-wake-dispatcher</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd /Users/tossaki/App/QuantRabbit-live &amp;&amp; PYTHONPATH=src /opt/homebrew/bin/python3 tools/guardian_wake_dispatcher.py</string>
  </array>
  <key>StartInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>/Users/tossaki/App/QuantRabbit-live/logs/guardian_wake_dispatcher.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/tossaki/App/QuantRabbit-live/logs/guardian_wake_dispatcher.launchd.err</string>
  <key>WorkingDirectory</key>
  <string>/Users/tossaki/App/QuantRabbit-live</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>src</string>
    <key>QR_GUARDIAN_WAKE_GATEWAY_HANDOFF</key>
    <string>0</string>
    <key>QR_GUARDIAN_ACTION_EXECUTE</key>
    <string>0</string>
    <key>QR_LIVE_ENABLED</key>
    <string>0</string>
    <key>AI_ORDER_AUTHORITY</key>
    <string>NONE</string>
    <key>QR_GUARDIAN_WAKE_SUPERVISOR_ONLY</key>
    <string>1</string>
    <key>QR_GUARDIAN_WAKE_TUNING_MODE</key>
    <string>HOURLY</string>
    <key>QR_GUARDIAN_WAKE_CODEX_BIN</key>
    <string>/Applications/ChatGPT.app/Contents/Resources/codex</string>
    <key>QR_GUARDIAN_WAKE_CODEX_PREFLIGHT</key>
    <string>1</string>
    <key>CODEX_DISABLE_UPDATE_CHECK</key>
    <string>1</string>
  </dict>
</dict>
</plist>
